# OS_fenikkusu_files/etc/systemd/system/save_data.service
[Unit]
Description=Sauvegarde

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/etc/fenikkusu/save_data

[Install]
WantedBy=multi-user.target