# Tatsuki-OS/system-files/etc/systemd/system/load_data.service
[Unit]
Description=Chargement

[Service]
Type=oneshot
ExecStart=/etc/usagi/load_data
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target