# Usagi-OS/system-files/etc/systemd/system/save_data.service
[Unit]
Description=Sauvegarde

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/etc/usagi/save_data

[Install]
WantedBy=multi-user.target